# livegrid/srx_livegrid_worker.conf
[program: srx_livegrid_worker]
command=/opt/conda_envs/bluesky_workers/bin/python /home/bluesky_workers/nslsii/srx-workers/livegrid/srx_livegrid_worker.py
user=bluesky_workers